# deploy/systemd/pdg-tproxy.service
[Unit]
Description=PrivDNS Gateway tproxy policy-routing + nftables
After=network-online.target sing-box.service
Wants=network-online.target
PartOf=sing-box.service

[Service]
Type=oneshot
RemainAfterExit=yes
# 策略路由: 打了 fwmark 的包查表 100, 走本地 (交给 sing-box tproxy)。
# 前置 del 用 '-' 忽略首次不存在的错误, 保证可重复执行。
ExecStartPre=-/usr/sbin/ip rule del fwmark 0x1 lookup 100
ExecStartPre=-/usr/sbin/ip -6 rule del fwmark 0x1 lookup 100
ExecStart=/usr/sbin/ip rule add fwmark 0x1 lookup 100
ExecStart=/usr/sbin/ip route replace local 0.0.0.0/0 dev lo table 100
ExecStart=/usr/sbin/ip -6 rule add fwmark 0x1 lookup 100
ExecStart=/usr/sbin/ip -6 route replace local ::/0 dev lo table 100
ExecStart=/usr/sbin/nft -f /var/lib/pdg/out/pdg.nft
ExecStop=-/usr/sbin/nft delete table inet pdg
ExecStop=-/usr/sbin/ip rule del fwmark 0x1 lookup 100
ExecStop=-/usr/sbin/ip -6 rule del fwmark 0x1 lookup 100

[Install]
WantedBy=multi-user.target
