# etc/systemd/5gpn-intercept.service
# 5gpn-unit-id: 5gpn-intercept.service:v1
[Unit]
Description=5gpn modular allowlisted MITM sidecar
After=network-online.target 5gpn-intercept-cert.service
Wants=network-online.target
Requires=5gpn-intercept-cert.service

[Service]
Type=simple
ExecCondition=/opt/5gpn/bin/5gpn-intercept --config /etc/5gpn/intercept/config.json --check-enabled
ExecStart=/opt/5gpn/bin/5gpn-intercept --config /etc/5gpn/intercept/config.json
Restart=on-failure
RestartSec=5
User=gpn-intercept
Group=gpn-intercept
UMask=0077
StateDirectory=5gpn-intercept
StateDirectoryMode=0700
LimitNOFILE=32768
MemoryHigh=512M
MemoryMax=768M
TasksMax=1024
NoNewPrivileges=yes
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
ProtectProc=invisible
ProcSubset=pid
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_INET AF_UNIX
ReadOnlyPaths=/etc/5gpn/intercept
InaccessiblePaths=-/etc/5gpn/intercept-ca -/etc/5gpn/acme -/etc/5gpn/dns.env -/etc/5gpn/mihomo

[Install]
WantedBy=multi-user.target
