# etc/systemd/5gpn-intercept-cert.path
# 5gpn-unit-id: 5gpn-intercept-cert.path:v1
[Unit]
Description=Watch the interception module snapshot for certificate changes

[Path]
PathChanged=/etc/5gpn/intercept/config.json
Unit=5gpn-intercept-cert.service

[Install]
WantedBy=multi-user.target
