# etc/systemd/mihomo.service
# 5gpn-unit-id: mihomo.service:v1
[Unit]
Description=5gpn mihomo (TLS external-controller + SNI/QUIC forward + egress)
After=network-online.target 5gpn-intercept.service
Wants=network-online.target 5gpn-intercept.service

[Service]
Type=simple
ExecStart=/opt/5gpn/bin/mihomo -f /etc/5gpn/mihomo/config.yaml -d /etc/5gpn/mihomo
Restart=on-failure
RestartSec=5
User=mihomo
Group=mihomo
UMask=0027
LimitNOFILE=65535
# Sandbox: leaf proxy — binds 80/443/8080/8443 tcp + 443 udp by default;
# optional catalog modules may add 5060 tcp/udp, plus loopback TLS external-controller
# on the shared zash role certificate. mihomo dials both IPv4 and IPv6;
# AF_NETLINK is REQUIRED — its UDP/QUIC DIRECT
# dial path does a route-table lookup (route ip+net: netlinkrib) that fatals the
# forward with "address family not supported by protocol" without it (confirmed
# on test-env 2026-07-14: QUIC/HTTP3 forwarding fails without AF_NETLINK; TCP is
# unaffected).
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
# mihomo v1.19.28 safe-path checks only trust files under -d plus SAFE_PATHS.
# The shared zash controller cert lives outside /etc/5gpn/mihomo, so allow only
# that read-only cert dir while keeping ProtectSystem=strict and writes scoped to
# /etc/5gpn/mihomo.
Environment=SAFE_PATHS=/etc/5gpn/cert/zash
# mihomo writes provider caches + reads whitelist.txt/config.yaml under its dir.
ReadWritePaths=/etc/5gpn/mihomo
# ACME credentials are unrelated to the data plane and must remain unreadable.
InaccessiblePaths=-/etc/5gpn/acme -/etc/5gpn/dns.env

[Install]
WantedBy=multi-user.target
