# deploy/ios/pdg-ios-profile-cleanup.timer
[Unit]
Description=PrivDNS Gateway expired iOS profile cleanup timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=pdg-ios-profile-cleanup.service

[Install]
WantedBy=timers.target
