# deploy/ios/pdg-ios-profile-sync.service
[Unit]
Description=PrivDNS Gateway sync registered client IPs to nftables
After=network-online.target nftables.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/pdg-bot/profile-http.py --sync
Environment=WWW_DIR=/opt/pdg-bot/ios-www
Environment=ACCESS_STATE_FILE=/opt/pdg-bot/ios-www/.ios-access.json
EnvironmentFile=-/etc/privdns-gateway/ios-profile.env
User=root
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/opt/pdg-bot/ios-www
RestrictAddressFamilies=AF_NETLINK AF_UNIX
RestrictSUIDSGID=yes
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
