# deploy/ios/pdg-probe81.service
[Unit]
Description=PrivDNS Gateway iOS OnDemand 探测端点 (:81 -> 200)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 /opt/pdg-bot/probe81.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
DynamicUser=true
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
